# vagrant/app/mailhog.service
[Unit]
Description=MailHog

[Service]
ExecStart=/bin/mailhog

[Install]
WantedBy=multi-user.target
